# bmm.service
[Unit]
Description=Gunicorn instance to serve BMM
After=network.target

[Service]
User=michaeliyke
Group=www-data
WorkingDirectory=/home/michaeliyke/michaeliyke/projects/BMM
Environment=PATH=/home/michaeliyke/michaeliyke/projects/BMM/bin
EnvironmentFile=/home/michaeliyke/michaeliyke/projects/BMM/.env
ExecStart=/home/michaeliyke/michaeliyke/projects/BMM/start.sh

[Install]
WantedBy=multi-user.target

